# upstart.conf
start on runlevel [234]
respawn
env RUST_LOG=warn
exec /usr/bin/lithos_tree

