# upstart.conf
start on runlevel [234]
respawn
env RUST_LOG=warn
env RUST_BACKTRACE=1
exec /usr/bin/lithos_tree

